# systemd/mcp-browser.socket
[Unit]
Description=MCP Browser Socket

[Socket]
ListenStream=/run/mcp-browser/mcp-browser.sock
SocketMode=0666

[Install]
WantedBy=sockets.target